# systemd/system/speedtest-logger.service
[Unit]
Description=Speedtest Telemetry Logger (single run)

[Service]
Type=oneshot
ExecStart=/app/scripts/speedtest-log.sh
WorkingDirectory=/app
