# systemd/user/speedtest-logger.service
[Unit]
Description=Speedtest Telemetry Logger (single run)

[Service]
Type=oneshot
ExecStart=%h/speedtest-logs/scripts/speedtest-log.sh
WorkingDirectory=%h/speedtest-logs
